# scripts/flyers/com.spmkt.classify-new.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.spmkt.classify-new</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/diogolopes/.local/share/fnm/node-versions/v24.12.0/installation/bin/node</string>
    <string>/Users/diogolopes/dev/Spmkt/scripts/flyers/classify-meat.mjs</string>
    <string>--new</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/diogolopes/dev/Spmkt</string>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/diogolopes/Library/Logs/spmkt-classify-new.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/diogolopes/Library/Logs/spmkt-classify-new.log</string>
</dict>
</plist>
